# system/etc/systemd/system/dnsmasq-restart.path
[Path]
PathChanged=/etc/hosts.coldturkey

[Install]
WantedBy=multi-user.target
